# modules/luci-mod-network/root/etc/uci-defaults/50_luci-mod-admin-full
#!/bin/sh

if [ "$(uci -q get luci.diag)" != "internal" ]; then
	host=""

	if [ -s /etc/os-release ]; then
		. /etc/os-release
		host="${HOME_URL:-${BUG_URL:-$OPENWRT_DEVICE_MANUFACTURER_URL}}"
		host="${host#*://}"
		host="${host%%/*}"
	fi

	uci -q batch <<-EOF >/dev/null
		set luci.diag=internal
		set luci.diag.dns='${host:-google.com.vn}'
		set luci.diag.ping='${host:-google.com.vn}'
		set luci.diag.route='${host:-google.com.vn}'
		commit luci
	EOF
fi

exit 0
